# gapd7300h/telnet_overlay_en_for_hw/telnetd
#!/bin/sh /etc/rc.common

START=94

start() {
        /usr/sbin/telnetd
}

stop() {
        killall telnetd
}

restart() {
        killall telnetd
        sleep 1
        /usr/sbin/telnetd
}

