# gapd7300h/davo/feeds/system/dvlog/files/dvlog.init
#!/bin/sh /etc/rc.common

START=12

start() {
	/usr/sbin/dvlog -B -P /var/run/dvlog.pid
}

stop() {
	kill $(cat /var/run/dvlog.pid) &> /dev/null
}

restart() {
	kill $(cat /var/run/dvlog.pid) &> /dev/null
	sleep 1
	/usr/sbin/dvlog -B -P /var/run/dvlog.pid
}

