# p201/boards/rtl8197F/etc.P201/inittab
# Boot-time system configuration/initialization script.
::sysinit:/etc/init.d/rcS

# Start an "askfirst" shell on the console (whatever that may be)
::respawn:-/bin/sh

::shutdown:/bin/dvflag -s 4 REBOOTING 1
